# packaging/linux/ryujin-doom.service
[Unit]
Description=DOOM for the ASUS ROG Ryujin III Extreme LCD
After=systemd-udevd.service
ConditionPathExists=/opt/ryujin-doom/ryujin-doom
ConditionPathExists=/var/lib/ryujin-doom/IWAD.WAD

[Service]
Type=simple
User=ryujin-doom
Group=ryujin-doom
WorkingDirectory=/var/lib/ryujin-doom
ExecStart=/opt/ryujin-doom/ryujin-doom -iwad /var/lib/ryujin-doom/IWAD.WAD
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
KillSignal=SIGTERM
UMask=0077

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
ReadWritePaths=/var/lib/ryujin-doom

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
